# scripts/com.esther.cursor-chat-export.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.esther.cursor-chat-export</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/esther/Code/cursor-chat/scripts/export_cursor_chats.py</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>23</integer>
    <key>Minute</key><integer>55</integer>
  </dict>
  <key>StandardOutPath</key><string>/Users/esther/Library/Logs/cursor-chat-export.out.log</string>
  <key>StandardErrorPath</key><string>/Users/esther/Library/Logs/cursor-chat-export.err.log</string>
  <key>RunAtLoad</key><true/>
</dict>
</plist>


